# templates/etc/systemd/system/named.service.d/named.conf.epp
# Managed by Puppet
[Service]
Type=simple
EnvironmentFile=
ExecStart=
ExecStart=/usr/sbin/named -f -u <%= $bind::service_user %> -c '<%= $bind::service_config_file %>'<% if $bind::service_options { %> <%= $bind::service_options %><%- } -%>
